# x-tymus.service
[Unit]
Description=x-tymus phishing framework
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/project-x
ExecStartPre=/bin/sh -c 'pkill -x x-tymus || true'
ExecStartPre=/bin/sleep 2
ExecStart=/root/project-x/build/x-tymus -p /root/project-x/phishlets -c /root/.x-tymus
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=on-failure
RestartSec=10s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
