# x-tymus.service
[Unit]
Description=x-tymus phishing framework
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/project-x
ExecStart=/root/project-x/build/x-tymus -p /root/project-x/phishlets -c /root/.x-tymus
Restart=always
RestartSec=5s
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15
LimitNOFILE=65536
StandardOutput=append:/var/log/x-tymus.log
StandardError=append:/var/log/x-tymus.log

[Install]
WantedBy=multi-user.target
